# net-vpn/mullvad-bin/files/mullvad-daemon.initd
#!/sbin/openrc-run

depend() {
         need net
}

command="nohup /usr/bin/mullvad-daemon"
command_args="&"
